# asterisk/rootfs/etc/services.d/web/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Web UI..."

exec python3 /app/server.py
