# asterisk/rootfs/etc/services.d/asterisk/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Asterisk service..."

exec asterisk -f -vvv
